# service/etc/systemd/system/cadsled.socket
[Unit]
Description=Socket to communicate with Cads Led

[Socket]
ListenStream=/run/cadsled.socket
