# cads-edge/service/etc/systemd/system/wireguardup.timer
[Unit]
Description=Indicates Cads connection to MHY server 

[Timer]
OnBootSec=10sec
OnUnitActiveSec=10sec
RemainAfterElapse=no